# scripts/backup_daily.cron
# Run daily backup at 02:30
30 2 * * * cd <your-install-dir>/knowledge-base-system && ./scripts/backup_create.sh >> <your-install-dir>/knowledge-base-system/backups/backup.log 2>&1
